# frpc/rootfs/etc/services.d/frpc/finish
#!/usr/bin/execlineb -S0
s6-svwait -d -t 5000 /var/run/s6/services/frpc
