# frpc/rootfs/etc/services.d/frpc/run
#!/usr/bin/with-contenv bash
set -e
exec /usr/bin/run.sh
